# deploy/systemd/autobot-research-microstructure.timer
[Unit]
Description=Capture AUTOBOT public microstructure research evidence every 15 minutes

[Timer]
OnCalendar=*-*-* *:07,22,37,52:00 UTC
Persistent=true
RandomizedDelaySec=45s
AccuracySec=30s
Unit=autobot-research-microstructure.service

[Install]
WantedBy=timers.target
