# deploy/systemd/autobot-research-microstructure.service
[Unit]
Description=AUTOBOT forward public microstructure research collection
Documentation=file:///opt/Projet_AUTOBOT/src/autobot/v2/research/forward_microstructure_collection.py
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service
ConditionPathExists=/opt/Projet_AUTOBOT/deploy/systemd/run-autobot-research-microstructure-collection.sh

[Service]
Type=oneshot
WorkingDirectory=/opt/Projet_AUTOBOT
Environment=AUTOBOT_REPO_DIR=/opt/Projet_AUTOBOT
ExecStart=/opt/Projet_AUTOBOT/deploy/systemd/run-autobot-research-microstructure-collection.sh
TimeoutStartSec=6min
Nice=10
IOSchedulingClass=idle
IOSchedulingPriority=7
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
