# deploy/systemd/autobot-research-data.service
[Unit]
Description=AUTOBOT public research market-data collection
Documentation=file:///opt/Projet_AUTOBOT/config/research_data_collection.yaml
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service
ConditionPathExists=/opt/Projet_AUTOBOT/config/research_data_collection.yaml

[Service]
Type=oneshot
WorkingDirectory=/opt/Projet_AUTOBOT
Environment=AUTOBOT_REPO_DIR=/opt/Projet_AUTOBOT
ExecStart=/opt/Projet_AUTOBOT/deploy/systemd/run-autobot-research-collection.sh
TimeoutStartSec=90min
# ``docker run`` is a daemon client: terminating only its parent shell can
# otherwise leave the isolated collector alive.  Stop and remove *only*
# containers explicitly labelled as this daily research job; the AUTOBOT
# runtime service, derivatives jobs and Docker Compose containers are never
# selected by this cleanup.
KillMode=control-group
TimeoutStopSec=45s
ExecStop=/bin/sh -c '/usr/bin/docker ps -q --filter label=autobot.job=research-daily | xargs -r /usr/bin/docker stop --timeout 20'
ExecStopPost=/bin/sh -c '/usr/bin/docker ps -aq --filter label=autobot.job=research-daily | xargs -r /usr/bin/docker rm -f'
Nice=10
IOSchedulingClass=idle
IOSchedulingPriority=7

# The service starts only an isolated public-data container. The trading
# container, runtime database, .env, and Docker Compose service are untouched.
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
